# rootfs/etc/services.d/nextdns/run
#!/command/with-contenv sh
# NextDNS service - runs as non-root user
# Logs are handled by s6-log via the log/ subdirectory

NEXTDNS_VERSION=$(nextdns version 2>/dev/null | awk '{print $NF}' || echo "unknown")
echo "Starting NextDNS version ${NEXTDNS_VERSION}..." >&2

exec s6-setuidgid dnsmasq nextdns run \
    ${NEXTDNS_ARGUMENTS} \
    -log-queries \
    -config ${NEXTDNS_ID} 2>&1
